# antidetect-backend.service
[Unit]
Description=Anti-Detect Browser Backend API
After=network.target postgresql.service redis.service docker.service
Wants=postgresql.service redis.service docker.service

[Service]
Type=simple
User=berellian
Group=berellian
WorkingDirectory=/home/berellian/antol/anti-detect-mvp/backend
Environment="DATABASE_URL=postgres://antidetect_user:ABCDqwer1234@localhost:5432/antidetect_db"
Environment="REDIS_URL=redis://localhost:6379"
Environment="JWT_SECRET=your-secret-key-change-in-production"
Environment="RUST_LOG=info"
ExecStart=/home/berellian/antol/anti-detect-mvp/backend/target/release/backend
Restart=always
RestartSec=10
StandardOutput=append:/home/berellian/antol/anti-detect-mvp/backend/backend.log
StandardError=append:/home/berellian/antol/anti-detect-mvp/backend/backend.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/berellian/antol/anti-detect-mvp/backend

[Install]
WantedBy=multi-user.target
